# conteneurs/nginx/conf/etc/services.d/nginx/run
#!/usr/bin/execlineb -P

foreground {
	nginx -g "daemon off;"
}
