# conteneurs/nginx/conf/etc/services.d/sshd/run
#!/bin/execlineb -P

foreground {
	/usr/sbin/sshd -D -e
}
